# scripts/fupan.crontab
# fupan 复盘自动发布兜底 —— 每 10 分钟扫描 /fupan、/jieli、/jingjia 的新报告，
# 落盘到 data/reviews/（本机即线上，落盘即生效）并提交+推送到 dev。
# 安装：crontab /home/ubuntu/fupan/scripts/fupan.crontab
# 查看日志：tail -f /home/ubuntu/fupan/.publish_reviews.log
#
# ⚠️ 2026-06-12 起的分工：agent 的定时运行（9:25 竞价、17:15 复盘）归 OpenClaw cron
# （openclaw cron list 查看），不要在这里加 agent 调度——会和 OpenClaw 双跑打架。
# 本文件只保留发布兜底：PostToolUse hook 只覆盖 Write/Edit 落 ~/claudeCode/*.md 的场景，
# bash 落盘、jieli 的 .cache 产物、hook 失败等盲区靠这条 10 分钟轮询托底。
PATH=/usr/local/bin:/usr/bin:/bin
*/10 * * * * /usr/bin/flock -n /tmp/publish_reviews.lock /usr/bin/python3 /home/ubuntu/fupan/scripts/publish_reviews.py >> /home/ubuntu/fupan/.publish_reviews.log 2>&1
